# etc/supervisor/gunicorn.ini
[program:gunicorn]
command=%(ENV_HOME)s/Envs/factoring/bin/gunicorn factoring.wsgi:application -w 5 -k gaiohttp -n gunicorn_aiohttp -b 0.0.0.0:8000
directory=%(ENV_HOME)s/Git/factoring
user=root
autostart=true
autorestart=true
redirect_stderr=true
