# etc/supervisor/gunicorn.ini
[program:gunicorn]
command=/home/ubuntu/Envs/keithxiaoy/bin/gunicorn main:app -k aiohttp.worker.GunicornUVLoopWebWorker -n gunicorn_aiohttp_uvloop -w 9 -b 0.0.0.0:8000
directory=/home/ubuntu/Git/keithxiaoy
user=root
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
killasgroup=true

stdout_logfile_maxbytes=200MB
stdout_logfile_backups=10
stdout_logfile =/var/log/supervisor/keithxxiaoy-supervisor.log